# guacd/guacd.conf
[program:guacd]
command 			= guacd -l 4822 -b 0.0.0.0 -f
priority 			= 400
autorestart 	    = true
